# supervisor/mysql.conf
[program:mysql]
command=/usr/sbin/service mysql start